# ops/systemd/asus-codex-container.service
[Unit]
Description=ASUS VSE-01 (Containerized Codex Engineering Node)
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/t79/ASUS/codex
ExecStart=/usr/bin/docker compose -f /home/t79/ASUS/codex/docker/docker-compose.yml up -d --build
ExecStop=/usr/bin/docker compose -f /home/t79/ASUS/codex/docker/docker-compose.yml down
TimeoutStartSec=0
User=t79
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
